# units/ipcd.service.in
[Unit]
Description=IPC Daemon
After=network.target

[Service]
ExecStart=@sbindir@/ipcd
Restart=on-failure
