# units/ipcd.service.in
[Unit]
Description=IPC Daemon
After=network.target

[Service]
ExecStart=@sbindir@/ipcd
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=ipcd.socket
