# units/ipcd.service.in
[Unit]
Description=IPC Daemon

[Service]
ExecStart=@sbindir@/ipcd
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
