# units/ipcd.socket
[Unit]
Description=IPC daemon socket

[Socket]
ListenStream=0.0.0.0:9443

[Install]
WantedBy=sockets.target
