# units/ipcd.socket
[Unit]
Description=IPC daemon socket
Before=ipcd.service

[Socket]
ListenStream=9443
Accept=yes
